# core/ping6.timer
[Unit]
Description=Verify IPv6

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=default.target
