# onfailure/failure-email@.service
[Unit]
Description=OnFailure for %i

[Service]
Type=oneshot
ExecStart=/home/user/.config/systemd/onfailure.sh %i
