# onfailure/failure-email@.service
[Unit]
Description=OnFailure for %i

[Service]
Type=oneshot
ExecStart=/usr/bin/env %h/.config/systemd/utils/onfailure/onfailure.sh %i
