# units/backup-daily.target
[Unit]
Description=Daily Backups
OnFailure=failure-email@%i.service
