# units/backup-daily.timer
[Unit]
Description=Daily Backups

[Timer]
OnCalendar=*-*-* 3:21:00
Persistent=true
Unit=backup-daily.target

[Install]
WantedBy=default.target
