# units/backup-daily.target
[Unit]
Description=Daily Backup Target
OnFailure=failure-email@%i.service
