# units/paccache.timer
[Unit]
Description=Clean-up Old Pacman Files

[Timer]
OnCalendar=daily
AccuracySec=1h

[Install]
WantedBy=multi-user.target
