# units/tahoe-check.timer
[Unit]
Description=Tahoe Deep Check and Repair

[Timer]
OnCalendar=Mon 3:29:00
Persistent=true

[Install]
WantedBy=default.target
